# package/liteon-apps/preip/files/preipd.init
#!/bin/sh /etc/rc.common
START=100
STOP=100

DAEMON_NAME=preip
boot() {
	start
}

start() {
	[ -e /usr/sbin/preip ] && {
		echo "$DAEMON_NAME start..."
		/usr/sbin/preip & 2>&1
	}
        
}

stop() {
	echo "$DAEMON_NAME stop..."
	killall preip
}
