# package/liteon-apps/preip/files/preipd.init
#!/bin/sh /etc/rc.common
START=100
STOP=100

DAEMON_NAME=preip
boot() {
	start
}

start() {
	[ -e /usr/sbin/preipd ] && {
		echo "$DAEMON_NAME start..."
		insmod /usr/sbin/preip.ko
		/usr/sbin/preipd & 2>&1
	}
        
}

stop() {
	echo "$DAEMON_NAME stop..."
	killall preipd
	rmmod preip
}
